# buildfiles/app.service
[Unit]
Description=app Service

[Service]
WorkingDirectory=/tmp
ExecStart=python -m app
